# .config/systemd/user/hyprsunset.service
[Unit]
Description=Run hyprsunset 

[Service]
Type=oneshot
ExecStart=/bin/bash /home/caio/.config/scripts/hyprsunset/hyprsunset_autostart.sh
